# me.jessedyck.restic-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>me.jessedyck.restic-backup</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/node</string>
            <string>~/restic-backup/restic-backup.js</string>
        </array>
        <key>StartInterval</key>
        <integer>3600</integer>
        <key>StandardOutPath</key>
        <string>~/me.jessedyck.restic-backup.log</string>
        <key>StandardErrorPath</key>
        <string>~/me.jessedyck.restic-backup.log</string>
    </dict>
</plist>